# tools/com.webkitium.ds4-runtime.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.webkitium.ds4-runtime</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>python3</string>
    <string>/Users/mac/webkitium/tools/ds4_runtime_supervisor.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/mac/webkitium</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/mac/webkitium/.runtime/ds4-supervisor/launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/mac/webkitium/.runtime/ds4-supervisor/launchd.err.log</string>
</dict>
</plist>
